# container/supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=10MB
logfile_backups=1

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=10

[program:backend]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend-stdout.log
stderr_logfile=/var/log/supervisor/backend-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=20
environment=
    HOMEPILOT_MODE="container",
    HOMEPILOT_LLM_BACKEND="%(ENV_HOMEPILOT_LLM_BACKEND)s",
    COMFY_BASE_URL="http://127.0.0.1:8188"

# ComfyUI — image + video diffusion workflow engine.
# Runs in its own venv (/opt/comfy-venv) with CPU-only torch wheels
# installed in stage 2 of the Dockerfile. Binds to loopback so only
# the backend can reach it; nginx does NOT expose :8188 to the host.
# Takes 15-30 s to warm up on first start while it loads custom
# nodes — the backend /api/health endpoint is served independently
# so the desktop splash does not wait for it.
[program:comfyui]
command=/opt/comfy-venv/bin/python /opt/ComfyUI/main.py --listen 127.0.0.1 --port 8188 --preview-method auto --disable-auto-launch
directory=/opt/ComfyUI
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/comfyui-stdout.log
stderr_logfile=/var/log/supervisor/comfyui-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=25
startsecs=10
startretries=3
