# container/supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=10MB
logfile_backups=1

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=10

[program:backend]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backend-stdout.log
stderr_logfile=/var/log/supervisor/backend-stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=20
environment=
    HOMEPILOT_MODE="container",
    HOMEPILOT_LLM_BACKEND="%(ENV_HOMEPILOT_LLM_BACKEND)s"
